# services.d/httpbin/run
#!/usr/bin/with-contenv sh

echo >&2 "starting httpbin"

exec /usr/bin/actix-httpbin 2>&1