# services.d/httpbox/run
#!/usr/bin/with-contenv sh

echo >&2 "starting httpbox"

exec /usr/bin/httpbox 2>&1